# deploy/reloop-rvm.service
[Unit]
Description=ReLoop RVM Edge Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=reloop
Group=reloop
WorkingDirectory=/opt/reloop
ExecStart=/opt/reloop/.venv/bin/python -m rvm.main --config /etc/reloop/rvm.json
Restart=always
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/reloop /var/log/reloop
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
